# default-conf/etc/crontab
*/15    * * * * www-data  (cd /home/engel-system/services ; php4 cron_dect.php)
