# default-conf/etc/crontab
*/15    * * * * www-data  (cd /var/www/services/ ; php5 cron_dect.php)
