# default-conf/etc/crontab
*/15    * * * * www-data  (cd /var/www.ssl/admin ; php4 dect_cron.php)
